# root/etc/uci-defaults/41_luci-cloudflared
#!/bin/sh


uci	delete ucitrack.@cloudflared[-1]
uci	add ucitrack cloudflared
uci	set ucitrack.@cloudflared[-1].init=cloudflared
uci	commit ucitrack



exit 0
